# etc/init.d/grabber-start
#!/sbin/runscript

description="Starting Grabber and respawn it if crashes."

depend() {
	after bootmisc
	need net
}

start() {
	ebegin "Starting Grabber launcher deamon"
	/root/grabber_start.sh &
	echo $! > /root/grabber_start.pid
	eend $?
}

stop() {
	ebegin "Stopping Grabber launcher daemon"
	kill -9 $(cat /root/grabber_start.pid)
	eend $?
}
